# config/dx-cluster-websocket.service
[Unit]
Description=DX Cluster Web WebSocket Server
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/om0rx/public_html/dx-cluster
ExecStart=/usr/bin/php /home/om0rx/public_html/dx-cluster/server/production-websocket-server.php
Restart=always
RestartSec=10

# Environment variables
Environment=PHP_ENV=production

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dx-cluster-websocket

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target